# roles/pve/files/bash/busybox.service
[Unit]
Description=My Script Service
After=network.target

[Service]
ExecStart=~/Sync/www/pve-bb.sh

[Install]
WantedBy=multi-user.target
